# pipeline/cron/com.id8labs.shipped-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Shipped. — Nightly status heartbeat
  Runs every night at 21:00 local. Regenerates the dashboard and posts
  a macOS notification with readiness %, gap count, article count.
  Does NOT publish or commit. Passive nudge only.

  Install via: bash pipeline/scripts/install-cron.sh
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.id8labs.shipped-nightly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/eddiebelaval/Development/id8/shipped/pipeline/scripts/nightly-check.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>21</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/eddiebelaval/Library/Logs/shipped/nightly-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/eddiebelaval/Library/Logs/shipped/nightly-stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/eddiebelaval/.nvm/versions/node/v22.21.1/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/eddiebelaval</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/eddiebelaval/Development/id8/shipped/pipeline</string>
</dict>
</plist>
